# config/includes.chroot/etc/systemd/system/umbra-services.service
[Unit]
Description=Umbra OS — ensure hardening services are active (AppArmor, resolved)
After=network.target apparmor.service systemd-resolved.service
Wants=apparmor.service systemd-resolved.service
# Run before the login manager so the live-vs-installed autologin fixup applies.
Before=display-manager.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Force AppArmor + a working resolver (falls back to plain DNS if the hardened
# path can't resolve — never leaves the box without DNS).
ExecStart=/usr/local/sbin/umbra-boot-net.sh

[Install]
WantedBy=multi-user.target
