# config/includes.chroot/etc/systemd/system/umbra-services.service
[Unit]
Description=Umbra OS — ensure hardening services are active (AppArmor, resolved)
After=network.target apparmor.service systemd-resolved.service
Wants=apparmor.service systemd-resolved.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Force-load AppArmor profiles and start encrypted-DNS resolver, even in the
# live environment where the normal enable-symlinks sometimes don't trigger.
ExecStart=/bin/sh -c 'systemctl restart apparmor.service 2>/dev/null || true; systemctl start systemd-resolved.service 2>/dev/null || true; ln -sf /run/systemd/resolve/stub-resolv.conf /etc/resolv.conf 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
